# templates/etc/init/logstash-forwarder.conf.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

description "Log shipper"
start on runlevel [2345]
stop on runlevel [!2345]

console log
limit nofile 65536 65536
respawn
setgid {{logstash_user}}
setuid {{logstash_group}}
chdir {{logstash_forwarder_home}}

exec {{logstash_forwarder_home}}/bin/logstash-forwarder -config {{logstash_confdir}}/forwarder.json
